# clearlinux/systemd/template-preparation.service
[Unit]
Description=Template preparation service

[Service]
ExecStart=/prepare-custom-template.sh

[Install]
WantedBy=multi-user.target

